# aquable/rootfs/etc/services.d/aquable/finish
#!/usr/bin/with-contenv bashio
# filepath: aquable/rootfs/etc/services.d/aquable/finish

bashio::log.info "AquaBle service stopped"
